# veyon-4.10.0/service/io.veyon.service.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.veyon.service</string>
    <key>ProgramArguments</key>
    <array>
        <string>@CMAKE_INSTALL_PREFIX@/bin/veyon-service</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/var/log/veyon-service.log</string>
    <key>StandardOutPath</key>
    <string>/var/log/veyon-service.log</string>
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
